# crontab
*/15 * * * /home/bmg_regex/btlscriptchallenge/checkfilescript.sh
